# deploy/ds4-server.service
[Unit]
Description=DS4 DSpark GB10 Server
Documentation=https://github.com/redhunt07/ds4-DSpark-GB10
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=redhunt07
Group=redhunt07
WorkingDirectory=/home/redhunt07/ai-stack/ds4-gb10-dspark

Environment=CUDA_VISIBLE_DEVICES=0
Environment=DS4_THREADS=10
# Accelerate only the N>1 F16 verifier GEMMs.  Unlike the legacy aggregate
# FAST_VERIFY switch, this leaves canonical N=1 decode and attention unchanged.
Environment=DS4_CUDA_FAST_VERIFY_F16=1
Environment=DS4_CUDA_MOE_NO_ATOMIC_DOWN=1
Environment=DS4_CUDA_NO_INDEXED_HEADS8=1
# Persistent keyed CUDA Graph for the fixed-shape N>1 verifier.  The runtime
# falls back to eager execution on capture/update failure.
Environment=DS4_GRAPH_MTP_VERIFY=1
# Inactivity watchdog for one single-flight request.  The deadline is renewed
# at every safe prefill frontier and decoded token, so a progressing long
# context is never stopped merely because of total wall-clock duration.
Environment=DS4_SERVER_JOB_TIMEOUT_SEC=3600
# Buffer tool deltas until DSML is complete so malformed calls can be retried
# without exposing irreversible partial JSON to OpenAI-compatible clients.
Environment=DS4_SERVER_BUFFER_TOOL_STREAM=1

MemoryHigh=118G
TasksMax=80
LimitNOFILE=1048576
IOSchedulingClass=best-effort
IOSchedulingPriority=2

ExecStart=/home/redhunt07/ai-stack/ds4-gb10-dspark/ds4-server \
  --model /home/redhunt07/ai-stack/ds4-gb10-dspark/gguf/DeepSeek-V4-Flash-DSpark-Abliterated-COMBINED-Q2-DS4.gguf \
  --cuda \
  --dspark \
  --ctx 131072 \
  --tokens 32768 \
  -t 10 \
  --prefill-chunk 1024 \
  --host 0.0.0.0 \
  --port 8000 \
  --cors \
  --warm-weights \
  --adaptive-thinking \
  --power 85 \
  --kv-disk-dir /home/redhunt07/ai-stack/ds4-gb10-dspark/kv-cache \
  --kv-disk-space-mb 98304 \
  --kv-cache-min-tokens 512 \
  --kv-cache-cold-max-tokens 120000 \
  --kv-cache-continued-interval-tokens 8192 \
  --kv-cache-boundary-trim-tokens 32 \
  --kv-cache-boundary-align-tokens 2048 \
  --tool-memory-max-ids 250000

Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ds4-dspark-gb10

[Install]
WantedBy=multi-user.target
